# contrib/croutond.service
[Unit]
Description=croutond VM orchestration daemon
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=croutond
Group=croutond
EnvironmentFile=-/etc/default/croutond
Environment=CLOUD_HYPERVISOR=/usr/bin/cloud-hypervisor
Environment=CROUTOND_RUNTIME_DIR=/var/lib/croutond
Environment=CROUTOND_POOL_SIZE=4
ExecStart=/usr/local/bin/croutond --ch-bin "${CLOUD_HYPERVISOR}" --runtime-dir "${CROUTOND_RUNTIME_DIR}" --pool-size "${CROUTOND_POOL_SIZE}"
Restart=on-failure
RestartSec=2
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=full
ReadWritePaths=/var/lib/croutond
StateDirectory=croutond
RuntimeDirectory=croutond
UMask=0077

[Install]
WantedBy=multi-user.target
